# com.droid.discord-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.droid.discord-bridge</string>
    
    <key>Program</key>
    <string>/Users/terry.gerighty/.factory/discord-bridge/bin/discord-bridge-daemon</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/discord-bridge.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/discord-bridge.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/terry.gerighty</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
